# memory-product/scripts/synthesis-cron.timer
[Unit]
Description=0Latency synthesis cron timer (daily 03:00 UTC)

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
